# deploy/aws/systemd/ai-incident-ui.service
[Unit]
Description=AI Incident Resolution UI (Streamlit)
After=network.target ai-incident-api.service
Requires=ai-incident-api.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/ai-incident-resolution
EnvironmentFile=/opt/ai-incident-resolution/.env
ExecStart=/opt/ai-incident-resolution/.venv/bin/streamlit run ui/app.py --server.port 8501 --server.address 127.0.0.1 --server.headless true --browser.gatherUsageStats false
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
